# files/supervisord-mysql.conf
[program:mysqld_safe]
command=/usr/bin/mysqld_safe --user=mysql --datadir="/data/db/mysql"
