# files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mysqld_safe]
command=/usr/bin/mysqld_safe --user=mysql --datadir="/data/db/mysql"
